# backend/deploy/arm-backend.service
[Unit]
Description=SC171V2 Arm Control FastAPI Backend
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hand-recognition/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/hand-recognition/backend/.env
ExecStart=/opt/hand-recognition/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
